# greenmind-gateway/systemd/greenmind-gateway.service
[Unit]
Description=GreenMind Raspberry Pi Edge Gateway
After=network-online.target NetworkManager.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# Root is required for nmcli AP management without Polkit configuration.
User=root
WorkingDirectory=/opt/greenmind/current
EnvironmentFile=-/opt/greenmind/current/.env
ExecStart=/opt/greenmind/current/venv/bin/python -m src.main

# Restart policy with crash-loop protection
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=greenmind-gateway

[Install]
WantedBy=multi-user.target
