# greenmind-gateway/systemd/greenmind-gateway.service
[Unit]
Description=GreenMind Raspberry Pi Edge Gateway
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
# User root is required to interact with nmcli (NetworkManager) for AP creation 
# without complex Polkit rules on a typical headless setup.
User=root 
WorkingDirectory=/opt/greenmind/gateway
EnvironmentFile=-/opt/greenmind/gateway/.env
ExecStart=/opt/greenmind/gateway/venv/bin/python src/main.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
